# templates/crontab
* * * * * ~/scripto-share/crontabs/ap_crontab_test.sh
